# rootfs/etc/profile.d/dlite.sh
#!/bin/sh

echo "dlite-os version $(cat /etc/version)"
docker -v
